# old/imagezmq.service
[Service]
ExecStart=/home/dek/src/microscope_ui/start_imagezmq.sh
User=dek
Before=graphical.target

[Install]
WantedBy=graphical.target
